# templates/kali-gui.service.j2
[Unit]
Description=Headless Sway + wayvnc + noVNC for Kali GUI
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User={{ kali_gui_user }}
WorkingDirectory=/home/{{ kali_gui_user }}

# systemd will create /run/kali-gui owned by {{ kali_gui_user }}
RuntimeDirectory=kali-gui
RuntimeDirectoryMode=0700
Environment=XDG_RUNTIME_DIR=/run/kali-gui

ExecStart=/usr/local/bin/kali-gui.sh
Restart=on-failure
RestartSec=2

[Install]
WantedBy=multi-user.target
